# crontab
# m h   dom mon dow   command
# create index file daily
0   7   *   *   *     touch /tmp/ping
